# shell/crontab/定时任务
#¶¨Ê±Ö´ÐÐ bakup.sh
0 0 * * * (sh /opt/shell/bakup.sh > /dev/null 2>&1)
